# apache_volume/etc/rc0.d/K01msmtpd
#!/usr/bin/env /lib/init/init-d-script
### BEGIN INIT INFO
# Provides:          msmtpd
# Required-Start:    $syslog $time
# Required-Stop:     $syslog $time
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: msmtp daemon
# Description:       Debian init script to start the msmtpd daemon
### END INIT INFO

INTERFACE=127.0.0.1
PORT=10025

NAME=msmtpd
START_ARGS="--background --make-pidfile"
DAEMON=/usr/bin/msmtpd
DAEMON_ARGS="--interface=${INTERFACE} --port=${PORT}"
